# crontab
* * * * * cd dev/rpi-home-display/ && ./update >> cron.log
